# ops/systemd/user/ops-snapshot-producer.service
[Unit]
Description=Generate the read-only Harness ops snapshot
ConditionFileIsExecutable=/home/slimy/slimy-harness/ops/snapshot-producer.sh

[Service]
Type=oneshot
WorkingDirectory=/home/slimy/slimy-harness
ExecStart=/home/slimy/slimy-harness/ops/snapshot-producer.sh
UMask=0027
NoNewPrivileges=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
