# ops/systemd/user/ops-snapshot-producer.timer
[Unit]
Description=Refresh the Harness ops snapshot every 10 minutes

[Timer]
OnActiveSec=10min
OnUnitActiveSec=10min
AccuracySec=30s
RandomizedDelaySec=15s
Unit=ops-snapshot-producer.service

[Install]
WantedBy=timers.target
